# systemd/user/oc-playwright-mcp.service
[Unit]
Description=Playwright MCP server (host-side, visible browser)
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.bin/oc-playwright-mcp
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
